# headphones-1.0.0/debian/headphones/DEBIAN/prerm
#!/bin/sh
# prerm script for headphones
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see /usr/share/doc/packaging-manual/

case "$1" in
	remove|upgrade|deconfigure)

		###
		# Disable our service
		###
		if [ `systemctl is-active fvapp-headphones` = "active" ]; then
			echo "stopping fvapp-headphones"
        		systemctl stop fvapp-headphones
		else 
        		echo "fvapp-headphones not running"
		fi

		if [ -f /lib/systemd/system/fvapp-headphones.service ]; then
			echo "disable fvapp-headphones"
        		systemctl disable fvapp-headphones
		else
        		echo "fvapp-headphones not enabled"
		fi

		###
		# Remove compiled python stuff
		###
		echo "remove compiled python stuff"
		find /apps/headphones -name "*.pyc" -print0 | xargs -n1 -0 rm -f
		find /apps/headphones -name "*.pyo" -print0 | xargs -n1 -0 rm -f
		echo "completed removing compiled python stuff"

		###
		# remove headphones app
		###
		if [ ! -f /apps/headphones ]; then
			echo "removing headphones app folder"
                	rm -rf /apps/headphones
			event_push app readynasd '<delete-s resource-type="LocalApp" resource-id="LocalApp"><LocalApp appname="headphones" success="1"/></delete-s>' 0 0
			echo "headphones Removed"
		else
			echo "Error headphones directory still present"
		fi
	;;
	failed-upgrade)
	;;
	*)
		echo "prerm called with unknown argument \`$1'" >&2
        	exit 0
    	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.


exit 0

